# overlay/etc/profile.d/trac.sh
#!/bin/bash

export PKG_RESOURCES_CACHE_ZIP_MANIFESTS=1
